# dot_config/systemd/user/random-wallpaper.service
[Service]
ExecStart=%h/.local/bin/random-wallpaper
